# discord-bot/cortex-discord-bot.service
[Unit]
Description=CortexPrism Discord Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/cortex-web/discord-bot
Environment=NODE_ENV=production
Environment=DOTENV_CONFIG_PATH=../.env
ExecStart=/root/.nvm/versions/node/v24.15.0/bin/node /root/cortex-web/node_modules/.bin/tsx src/index.ts
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
